# deploy/systemd/uv-im-connector.service
[Unit]
Description=uv-im-connector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=uvim
Group=uvim
EnvironmentFile=/etc/uv-im-connector/uv-im-connector.env
ExecStart=/usr/local/bin/uv-im-connector
Restart=always
RestartSec=5s
WorkingDirectory=/var/lib/uv-im-connector
StateDirectory=uv-im-connector
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/uv-im-connector

[Install]
WantedBy=multi-user.target
